# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:php5-fpm]
command=/usr/sbin/php5-fpm -c /etc/php5/fpm

[program:varnishd]
command=/usr/sbin/varnishd -f /etc/varnish/default.vcl -a 0.0.0.0:80 -s malloc,64m -F
autorestart=true
